# supervisord.conf
[supervisord]
nodaemon=true

[program:consul]
command=consul agent -server -bootstrap-expect 1 -data-dir /tmp/consul -node=agent-one